# installer/descriptor-poll.cron
* * * * * cd /home/descriptors && python3 /home/orchestrator.py
